# supervisord.conf
[supervisord]
nodaemon=true

[program:webhook]
command=/app/bin/webhook
autostart=true
autorestart=true
stderr_logfile=/app/logs/webhook.err.log
stdout_logfile=/app/logs/webhook.out.log
